# etc/khonliang-bus/khonliang-bus.service.d/webhook-secret.conf.example
# systemd drop-in: feed GITHUB_WEBHOOK_SECRET into the khonliang-bus unit.
#
# Install path:
#   /etc/systemd/system/khonliang-bus.service.d/webhook-secret.conf
#   owned by root:root, mode 0644
#
# After installing (and after editing /etc/khonliang/webhook-secret.env to
# carry the real secret), reload + restart:
#
#   sudo systemctl daemon-reload
#   sudo systemctl restart khonliang-bus.service
#
# Verify the bus picked the secret up — the webhook endpoint should now
# return 401 (signature missing) on an unsigned POST instead of 503
# (no secret configured). The port below assumes the canonical
# systemd-unit override of ``--port 8788``; the bare ``python -m bus``
# default is 8787, so substitute if your ExecStart doesn't pass --port:
#
#   curl -X POST http://localhost:8788/v1/webhooks/github \
#     -H 'X-GitHub-Event: ping' -d '{"zen":"smoke"}' \
#     -w '\nHTTP %{http_code}\n'

[Service]
EnvironmentFile=/etc/khonliang/webhook-secret.env
